# MCP-test/mcp-test.service
[Unit]
Description=Gunicorn instance to serve MCP-test Flask app
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/dataexploit-domain/MCP-test
Environment="PATH=/path/to/your/dataexploit-domain/MCP-test/venv/bin"
ExecStart=/path/to/your/dataexploit-domain/MCP-test/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5001 wsgi:app
Restart=always

[Install]
WantedBy=multi-user.target